# chef-repo/project_cookbooks/vx_solr/templates/default/initd.erb
#!/bin/bash

### BEGIN INIT INFO
# Provides:          solr
# Required-Start:
# Required-Stop:
# Default-Start:     3 5
# Default-Stop:      0 1 2 6
# Short-Description: init script for solr search
# Description:       init script for solr search
### END INIT INFO

# source function library
. /etc/rc.d/init.d/functions

LOG_FILE=<%= @log_file %>

NAME="Solr"

DESC="start/stop Solr Server"

start() {
  action "Starting ${NAME} instance 1: " <%= @script %> start \
  -cloud -p <%= @solr_port %> -z <%= @zk_config %> -s solr -V \
  <% unless @jvm_params.nil? then %> -a "<%= @jvm_params %>" <% end %> \
  <% unless @memory.nil? then %> -m "<%= @memory %>" <% end %>

  <% 2.upto(@solr_instances) do |instance| %>
    action "Starting ${NAME} instance <%= @instance %>: " <%= @script %> start \
    -cloud -p <%= @additional_instance_base_port + (instance - 2) %> -z <%= @zk_config %> -s home<%= instance %> -V \
    <% unless @jvm_params.nil? then %> -a "<%= @jvm_params %>" <% end %> \
    <% unless @memory.nil? then %> -m "<%= @memory %>" <% end %>
  <% end %>
}

stop() {
  action "Stopping ${NAME} instance 1: " <%= @script %> stop -p <%= @solr_port %>

  <% 2.upto(@solr_instances) do |instance| %>
    action "Stopping ${NAME} instance <%= @instance %>: " <%= @script %> stop -p <%= @additional_instance_base_port + (instance - 2) %>
  <% end %>
}

case $1 in
  start)
    start
    ;;
  stop)
    stop
    ;;
  restart)
    stop
    sleep 2
    start
    ;;
  status)
    <%= @script %> status | grep Found
    ;;
  *)
    echo "Usage: $0 {start|stop|restart}"
    exit 3
    ;;
esac
