# ehmp/product/production/jds/src/ewdjs/ewdrest.conf
description "EWD.js REST Service Upstart File"
author      "David Wicksell"

start on runlevel [2345]
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 10 5

# Change the value of USER below to the appropriate user that will run EWD.js

env USER="cacheusr"

script
    exec /bin/bash -c 'cd /home/$USER/ewdjs && node ewdRest.js' >> /var/log/ewdrest.log 2>&1
end script
